# templates/systemd/bento-sync.service
[Unit]
Description=Bento Tools Sync Daemon
Documentation=https://setup.divizend.com
After=network.target bento.service
Requires=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bento-sync
Environment="TOOLS_ROOT_GITHUB=${TOOLS_ROOT_GITHUB}"
Environment="BENTO_API_URL=http://localhost:4195"
Environment="S2_BASIN=${S2_BASIN}"
Environment="BASE_DOMAIN=${BASE_DOMAIN}"
Environment="S2_ACCESS_TOKEN=${S2_ACCESS_TOKEN}"
Environment="RESEND_API_KEY=${RESEND_API_KEY}"
ExecStart=/usr/local/bin/bun sync.ts
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
